# mwc/romana/source/4.2.x/usr/spool/cron/crontabs/uucp,v
head	4.2;
access
	hal
	ljg
	fred
	udo
	nigel;
symbols
	r12:4.2;
locks; strict;
comment	@# @;


4.2
date	92.07.01.00.18.16;	author rcsadmin;	state Exp;
branches;
next	;


desc
@@


4.2
log
@This is the baseline version
@
text
@# Cron Version 2.0. Installed on Tue Jun 30 14:49:24 1992
# Source was crontab.uucp

# Clean uucp logs at 3 AM weekdays.
00 03 * * 1-5	sh /usr/lib/uucp/uumvlog 5

# Conditionally poll site mwcbbs every day.
# Remove the comment delimiter and change the times on the next line as desired.
# 15 20 * * *	/usr/lib/uucp/uucico -smwcbbs

# end of uucp crontab
@
