# mwc/romana/relic/f/usr/spool/cron/crontabs/root
# Cron Version 2.0. Installed on Tue Jun 30 14:49:24 1992
# Source was crontab.root

# Execute atrun every five minutes.
0,5,10,15,20,25,30,35,40,45,50,55 * * * *	/usr/lib/atrun

# Rebuild the qfind database at 2 AM daily.
# Remove the comment delimiter on the command line below as desired.
00 02 * * *	/bin/qfind -b

# build 286 source archives on Wednesday night
#00 21 * * 3 /usr/bin/bld286

# build 386 source archives on Thursday night
#00 21 * * 4 /usr/bin/bld386

# end of root crontab
