# mwc/romana/relic/f/usr/spool/cron/crontabs/uucp
# Cron Version 2.0. Installed on Tue Jun 30 14:49:24 1992
# Source was crontab.uucp

# Clean uucp logs at 3 AM weekdays.
00 03 * * 1-5	sh /usr/lib/uucp/uumvlog 5

# Conditionally poll site mwcbbs every day.
# Remove the comment delimiter and change the times on the next line as desired.
# 15 20 * * *	/usr/lib/uucp/uucico -smwcbbs

# end of uucp crontab
